# script/baagh.service
[Unit]
Description= Baagh application

[Service]
ExecStart=/bin/baagh

[Install]
WantedBy=multi-user.target